# tools/systemd-units/v-services-restart.service
[Unit]
Description=V 长期服务的常驻守护 (Ollama / VCPToolBox / AgentTeam Board / agent-symphony / VCP admin)
# 14:21 根因修复: 用 fuguang user 跑 (user pip 路径)
# oneshot 改为 simple + Restart=always: 任何挂掉自动拉起
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=fuguang
Group=fuguang
WorkingDirectory=/home/fuguang
Environment="PATH=/home/fuguang/.local/bin:/home/fuguang/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/fuguang/.local/lib/python3.12/site-packages"
# 监控 + 拉起：用一个独立脚本跑, 比 oneshot shell 干净
ExecStart=/bin/bash /home/fuguang/.openclaw/workspace/tools/v-services-watchdog.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
