# tools/systemd-units/v-services-restart.service
[Unit]
Description=V 长期服务的轻量守护 (Ollama / VCPToolBox / AgentTeam Board)
# 不依赖 ollama/vcp/agentteam，只在 libvirtd 起来后兜底拉起
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# 检查 + 拉起三个服务（任一挂了就 nohup 重启）
ExecStart=/bin/bash -c '\
  set +e; \
  curl -sS -m 2 http://127.0.0.1:11434 >/dev/null 2>&1 || (cd /home/fuguang && setsid nohup /home/fuguang/bin/ollama serve >/home/fuguang/.ollama/server.log 2>&1 </dev/null & disown); \
  curl -sS -m 2 -o /dev/null -w "" http://127.0.0.1:6005 2>/dev/null; [ $? -ne 7 ] || (cd /home/fuguang/VCPToolBox && setsid nohup node server.js >/tmp/vcp.log 2>&1 </dev/null & disown); \
  curl -sS -m 2 http://127.0.0.1:8080 >/dev/null 2>&1 || (cd /home/fuguang/AgentTeam && setsid nohup /home/fuguang/.local/bin/agentteam board serve --port 8080 >/tmp/agentteam.log 2>&1 </dev/null & disown); \
  exit 0'
RemainAfterExit=yes
StandardOutput=journal

[Install]
WantedBy=multi-user.target
